# scripts/poolex.service
[Unit]
Description=Poolex Heat Pump RS485 Control
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/poolex-control
Environment="POOLEX_SERIAL_PORT=/dev/ttyUSB0"
Environment="POOLEX_DB_PATH=/var/lib/poolex/poolex.db"
Environment="POOLEX_API_PORT=5000"
Environment="POOLEX_MQTT_HOST=192.168.1.63"
Environment="POOLEX_MQTT_PORT=1883"
Environment="POOLEX_MQTT_USER=poolex"
Environment="POOLEX_MQTT_PREFIX=poolex"
EnvironmentFile=-/opt/poolex-control/poolex.env
ExecStart=/opt/poolex-control/venv/bin/python -m poolex
Restart=on-failure
RestartSec=15s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
