# scripts/poolex.service
[Unit]
Description=Poolex Heat Pump RS485 Control
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=/opt/poolex-control
Environment="POOLEX_SERIAL_PORT=/dev/ttyUSB0"
Environment="POOLEX_DB_PATH=/var/lib/poolex/poolex.db"
Environment="POOLEX_API_PORT=5000"
Environment="POOLEX_MQTT_HOST=localhost"
Environment="POOLEX_MQTT_PORT=1883"
# Environment="POOLEX_MQTT_USER=poolex"
# Environment="POOLEX_MQTT_PASSWORD=secret"
Environment="POOLEX_MQTT_PREFIX=poolex"
ExecStart=/opt/poolex-control/venv/bin/python -m poolex
Restart=on-failure
RestartSec=15s
StartLimitIntervalSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
